# dot_config/systemd/user/update-ssh-keys.timer
[Unit]
Description=Update SSH Keys every 10 minutes

[Timer]
OnBootSec=15s
OnUnitActive=10min

[Install]
WantedBy=timers.target

