# dot_config/systemd/user/update-ssh-keys.service
[Unit]
Description=Update User SSH Keys

[Service]
ExecStart=%h/.local/bin/update-ssh-keys acuteaura
Type=oneshot
